# wireguard/rootfs/etc/services.d/status/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard
# Shows current WireGuard status
# ==============================================================================
sleep 30
bashio::log.info "Requesting current status from WireGuard..."
exec wg show
